# .config/systemd/user/kak.service
[Unit]
Description=Start Kakoune in daemon mode

[Service]
Type=forking
ExecStart=/usr/bin/kak -s %u -d

[Install]
WantedBy=multi-user.target
